# tests/wpt/tests/tools/third_party/websockets/example/deployment/nginx/supervisord.conf
[supervisord]

[program:websockets-test]
command = python app.py
process_name = %(program_name)s_%(process_num)02d
numprocs = 4
autorestart = true
